# deploy/dms-webhook.service
[Unit]
Description=DMS Webhook Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/data/EquipmentManagement
ExecStart=/usr/bin/python3 /data/EquipmentManagement/scripts/webhook_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
